# src/system_systemd/user-services-setup.service
[Unit]
Description=Create/update user services from templates
After=network.target nginx.service

[Service]
Type=oneshot
ExecStart=/etc/user-services/scripts/create-user-services.sh
StandardOutput=journal
StandardError=journal
User=root
Group=root

# Sandboxing
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/etc/user-services /etc/nginx/user-services /etc/systemd/user
