# src/templates/code-server.socket.template
[Unit]
Description=Code Server Socket for %i
PartOf=code-server.service

[Socket]
ListenStream=/run/user/%UID%/code-server.sock
SocketUser=%i
SocketGroup=%i
SocketMode=0660
Accept=no
IdleTimeoutSec=3600

# Service to activate on connection
Service=code-server.service

# Keep socket alive even after service stops
KeepAlive=yes

[Install]
WantedBy=sockets.target
